# contrib/init/shibacoind.conf
description "Shibacoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SHIBACOIND_BIN="/usr/bin/shibacoind"
env SHIBACOIND_USER="shibacoin"
env SHIBACOIND_GROUP="shibacoin"
env SHIBACOIND_PIDDIR="/var/run/shibacoind"
# upstart can't handle variables constructed with other variables
env SHIBACOIND_PIDFILE="/var/run/shibacoind/shibacoind.pid"
env SHIBACOIND_CONFIGFILE="/etc/shibacoin/shibacoin.conf"
env SHIBACOIND_DATADIR="/var/lib/shibacoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # shibacoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SHIBACOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run shibacoind."
        echo "The setting must appear in $SHIBACOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Shibacoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SHIBACOIND_PIDDIR"
    chmod 0755 "$SHIBACOIND_PIDDIR"
    chown $SHIBACOIND_USER:$SHIBACOIND_GROUP "$SHIBACOIND_PIDDIR"
    chown $SHIBACOIND_USER:$SHIBACOIND_GROUP "$SHIBACOIND_CONFIGFILE"
    chmod 0660 "$SHIBACOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SHIBACOIND_PIDFILE" \
    --chuid $SHIBACOIND_USER:$SHIBACOIND_GROUP \
    --exec "$SHIBACOIND_BIN" \
    -- \
    -pid="$SHIBACOIND_PIDFILE" \
    -conf="$SHIBACOIND_CONFIGFILE" \
    -datadir="$SHIBACOIND_DATADIR" \
    -disablewallet \
    -daemon

